# units/user/tizen-user-middleware.target
[Unit]
Description=Tizen user middleware target
After=default.target
Requires=default.target
