# units/user/tizen-user-middleware-services.service
[Unit]
Description=Tizen user middleware services

[Service]
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/dbus/user_bus_socket
ExecStartPre=/usr/bin/sleep 0.5
ExecStart=/usr/bin/systemctl start --user tizen-user-middleware.target
