# units/user/tizen-user-middleware-services.path
[Unit]
Description=Path activation for tizen user middleware services

[Path]
PathExists=/tmp/fastboot

[Install]
WantedBy=default.target
