# units/system/tizen-middleware.target
[Unit]
Description=Tizen system middleware services
Conflicts=rescue.service rescue.target
Requires=multi-user.target
After=multi-user.target
